# com.dubsopenhub.ghost-ops.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key><string>com.dubsopenhub.ghost-ops</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/greggcochran/ghost-ops/ghost_ops.py</string>
    </array>
    <key>RunAtLoad</key><true/>
    <key>KeepAlive</key><true/>
    <key>StandardOutPath</key><string>/Users/greggcochran/ghost-ops/logs/ghost_ops.log</string>
    <key>StandardErrorPath</key><string>/Users/greggcochran/ghost-ops/logs/ghost_ops.err</string>
    <key>WorkingDirectory</key><string>/Users/greggcochran/ghost-ops</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key><string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
